# rc.local
#!/bin/sh
/home/mendel/clank/config.sh
su - mendel -c "/home/mendel/clank/start.sh"
